# supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=php-fpm

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"